# scripts/kicktipp.ini
[program:kicktipp]
directory=/home/kicktipp/kicktipp-deluxe
command=/usr/bin/node_modules/.bin/next start --port 3000
autostart=yes
autorestart=yes
startsecs=15
stderr_logfile=/home/kicktipp/logs/supervisord/kicktipp.err
stdout_logfile=/home/kicktipp/logs/supervisord/kicktipp.log
environment=NODE_ENV="production",PORT="3000"
